# crontab/periodical_crawl.sh
#!/bin/bash
# ./venv/bin/python ./main.py
#/home/ubuntu/digilog/digilog-crawling/crontab/venv/bin/python /home/ubuntu/digilog/digilog-crawling/crontab/main.py
/home/gerj/Documents/GitHub/digilog-crawling/crontab/venv/bin/python /home/gerj/Documents/GitHub/digilog-crawling/crontab/main.py